# systemd/gitlab.target
[Unit]
Description=GitLab

[Install]
WantedBy=multi-user.target
